# vm-systemd/qubes-power-reboot.service

[Unit]
Description=Reboot dom0 from within a VM
DefaultDependencies=no
Before=reboot.target

[Service]
Type=oneshot
ExecStart=/usr/bin/qrexec-client-vm dom0 qubes.Reboot
TimeoutStartSec=0

[Install]
WantedBy=reboot.target